# service/linux/socks5-proxy.service
[Unit]
Description=Socks 5 Proxy server
[Service]
ExecStart=/usr/local/bin/server
[Install]
WantedBy=multi-user.target
